# private_dot_config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Type=Application
Terminal=true
Comment=Run unclutter
Exec=unclutter -b --timeout 3
